# qm.container
[Service]
AllowedCPUs=6-11
CPUWeight=50
Delegate=true
IOWeight=50
ManagedOOMSwap=kill
MemorySwapMax=0
Slice=QM.slice
Restart=always
# Containers within the qm y default set OOMScoreAdj to 750
OOMScoreAdjust=500 

[Container]
AddCapability=all
ContainerName=qm
Exec=/sbin/init
Network=host
PodmanArgs=--security-opt label=nested --security-opt unmask=all
ReadOnly=true
Rootfs=/usr/lib/qm/rootfs
SecurityLabelFileType=qm_file_t
SecurityLabelType=qm_t
Volume=qmEtc:/etc:copy
Volume=qmVar:/var:copy
