# subsystems/wayland/etc/containers/systemd/qm-dbus-broker.container
# QM D-Bus Broker Container - Session bus for containerized applications
# Uses dbus-broker for improved performance over traditional dbus-daemon

[Unit]
Description=qm-dbus-broker container
After=qm-dbus.socket
Requires=qm-dbus.socket

[Container]
ContainerName=qm-dbus-broker
Image=quay.io/qm-images/qm-dbus-broker

# Launch dbus-broker in user session scope
Exec=/usr/bin/dbus-broker-launch --scope user

# D-Bus infrastructure mounts
Volume=/run/user/0/bus:/run/user/0/bus
Volume=/run/systemd:/run/systemd:ro
Volume=/run/user/0:/run/user/0
Volume=/etc/machine-id:/etc/machine-id:ro
Environment=XDG_RUNTIME_DIR=/run/user/0
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/0/bus
SecurityLabelType=qm_container_wayland_t

[Service]
Sockets=qm-dbus.socket
Restart=always

[Install]
Alias=qm-dbus.service
WantedBy=multi-user.target

