# subsystems/kvm/kvm.container
[Unit]
Description=kvm Container
After=network.target

[Container]
Image=quay.io/qm-images/kvm:latest

[Install]
WantedBy=multi-user.target
